# systemd/rw-metrics-export.service
[Unit]
Description=rw-backup-full: metrics export for Grafana
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/rw-backup-restore/scripts/metrics/metrics-exporter.sh
Nice=15
IOSchedulingClass=idle
TimeoutStartSec=30m
