# systemd/rw-sandbox-verify.service
[Unit]
OnFailure=rw-notify-failure@%n.service
Description=rw-backup-full: sandbox backup verification
After=docker.service network-online.target rw-backup-web.service
Requires=docker.service
Wants=network-online.target
# Не Requires=rw-backup-web.service: на хостах без компонента web юнита
# может не существовать вовсе — After= в этом случае просто не действует,
# а Requires= намертво отказал бы стартовать таймер.

[Service]
Type=oneshot
ExecStart=/opt/rw-backup-restore/scripts/sandbox/verify-entry.sh
Nice=15
IOSchedulingClass=idle
TimeoutStartSec=3h
